# library_masa.timer
[Unit]
Description=Library Service Timer for Masa

[Timer]
OnCalendar=*-*-* 2:30:00
Persistent=true

[Install]
WantedBy=timers.target